# systemd/serverpatcher.service
[Unit]
Description=Server Patcher - Linux patch automation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/serverpatcher run-once --config /etc/serverpatcher/config.json
# Hardening (review for your environment; patching requires broad system access)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes

[Install]
WantedBy=multi-user.target
